# sources/meta-fus/recipes-config/images/files/inittab
# /etc/inittab: init(8) configuration.
# $Id: inittab,v 1.91 2002/01/25 13:35:21 miquels Exp $

# The default runlevel.
id:5:initdefault:

# Boot-time system configuration/initialization script.
# This is run first except when booting in emergency (-b) mode.
si::sysinit:/etc/init.d/rcS

# What to do in single-user mode.
~~:S:wait:/sbin/sulogin

# /etc/init.d executes the S and K scripts upon change
# of runlevel.
#
# Runlevel 0 is halt.
# Runlevel 1 is single-user.
# Runlevels 2-5 are multi-user.
# Runlevel 6 is reboot.

l0:0:wait:/etc/init.d/rc 0
l1:1:wait:/etc/init.d/rc 1
l2:2:wait:/etc/init.d/rc 2
l3:3:wait:/etc/init.d/rc 3
l4:4:wait:/etc/init.d/rc 4
l5:5:wait:/etc/init.d/rc 5
l6:6:wait:/etc/init.d/rc 6
# Normally not reached, but fallthrough in case of emergency.
z6:6:respawn:/sbin/sulogin
# Put a getty on the login_tty. login_tty is a symbolic link to the
# device given in /sys/bdinfo/login_tty. This link is created by
# /etc/init.d/S01fssetup. /sbin/fsgetty is a script that simply starts
# the real /sbin/getty if the login_tty is a valid string. Otherwise it
# stays sleeping in an endless loop, because if it would return here, it
# would respawn immediately, and fail again and again, causing lots of
# annoying syslog messages.
_fus::respawn:/sbin/fsgetty
#mxc::respawn:/etc/init.d/rc_mxc.S
#mxc1:12345:respawn:/sbin/getty 115200 ttymxc1
# /sbin/getty invocations for the runlevels.
#
# The "id" field MUST be the same as the last
# characters of the device (after "tty").
#
# Format:
#  <id>:<runlevels>:<action>:<process>
#
# we at F&S assume that we have only one tty the login_tty.
#1:2345:respawn:/sbin/getty 38400 tty1

gpu::sysinit:/etc/init.d/rc_gpu.S
